# systemd/bird-updater.timer
[Unit]
Description=Bird Prefix Updater — daily 05:00
Requires=bird-updater.service

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true
AccuracySec=5min

[Install]
WantedBy=timers.target
